# whisper_dictation/ydotoold.service
[Unit]
Description=ydotoold service

[Service]
ExecStart=/usr/bin/ydotoold
Restart=always

[Install]
WantedBy=default.target